# contrib/upstart/upstart.conf
start on started rc
stop on stopped rc

respawn
respawn limit 3 5

setuid edms
setgid edms

exec /home/nucleos/production/deploy/production/service_demon.sh